# data/overlayfiles/systemd-sapinit/etc/systemd/system/sapinit.service.d/type.conf
# SAP still uses SysV init scripts.
# All SAP instances will be stopped if saphostexec is stopped or restarted.
# This file ensures the expected behavior

[Service]
Type=oneshot
